# plists/com.content-engine.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.content-engine.cleanup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nfainstein/.nvm/versions/node/v22.18.0/bin/node</string>
    <string>/Users/nfainstein/content-engine/dist/jobs/cleanup.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/nfainstein/content-engine</string>
  <key>StandardOutPath</key>
  <string>/Users/nfainstein/content-engine/logs/cleanup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nfainstein/content-engine/logs/cleanup.error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/nfainstein/.nvm/versions/node/v22.18.0/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
